# cronzimus-supervisord.ini
[program:cronzimus]
command=/usr/local/bin/gunicorn --bind 0.0.0.0:8080 -w 1 cronzimus.app:main() --access-logfile /usr/src/app/server.log
directory=/usr/src/app/
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisord/cronzimus-supervisord.log
stdout_logfile_maxbytes=1024MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisord/cronzimus-supervisord.log
stderr_logfile_maxbytes=1024MB
stderr_logfile_backups=5